# 1.8/s6/tail/run
#!/usr/bin/env bash

set -e

# Defaults
: ${LOGOUTPUT:='/dev/stdout'}

if [[ ${LOGOUTPUT} == "/dev/stdout" ]]; then
  s6-svc -d $(pwd)
  exit
else
  exec tail -F "${LOGOUTPUT}"
fi
